# packages/swww/systemd/swww-update.timer
[Unit]
Description=swww-update

[Timer]
OnCalendar=*-*-* *:00/5:00
Persistant=true

[Install]
WantedBy=timers.target